# root_copy/etc/cron.d/bdscore
0 */2 * * 0 root cd /opt/bds_core_vm && git pull && cp -rfv root_copy/* / && systemctl daemon-reload